# deploy/traxus-server.service
[Unit]
Description=Traxus WebSocket Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Traxus
ExecStart=/home/ubuntu/Traxus/.venv/bin/python -m server.main
Environment=TRAXUS_HOST=127.0.0.1
Environment=TRAXUS_PORT=8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
